# renaRatchet.conf
[program:renaRatchet]
command=/usr/bin/php /storage/www/projectRena/Rena run:ratchet
directory=/storage/www/projectRena
autostart=true
autorestart=true
user=www-data
group=www-data
