# renaCron.conf
[program:renaCron]
command=/usr/bin/php /storage/www/projectRena/Rena cron:run
directory=/storage/www/projectRena
autostart=true
autorestart=true
user=www-data
group=www-data
