# renaStompReceive.conf
[program:renaStomp]
command=/usr/bin/php /storage/www/projectRena/Rena stomp:run
directory=/storage/www/projectRena
autostart=true
autorestart=true
user=www-data
group=www-data

