# external/renaZKillboard.conf
[program:renaZKB]
command=/usr/bin/php /storage/www/projectRena/Rena run:zkb
directory=/storage/www/projectRena
autostart=true
autorestart=true
user=www-data
group=www-data
